# test/src/ubuntu/x86_64/.config/systemd/user/sysfact-apply.timer
[Unit]
Description=Periodically run Sysfact Apply

[Timer]
OnBootSec=30s
OnUnitActiveSec=10min

[Install]
WantedBy=timer.target
