# test/src/ubuntu/x86_64/.config/systemd/user/sysfact-apply.service
[Unit]
Description=Sysfact Apply

[Service]
Type=oneshot
ExecStart=/usr/bin/sysfact -L info apply

[Install]
WantedBy=default.target

